# crontab
@reboot sh /home/blue/sensors/temp_cuisine.sh >/home/blue/sensors/logs/cronlog 2>&1
@reboot sh /home/blue/sensors/hum_cuisine.sh >/home/blue/sensors/logs/cronlog 2>&1
